# config/home/.config/autostart/lxsession.desktop
[Desktop Entry]
Exec=/usr/bin/lxsession
Name=lxsession
Type=Application
