# config/home/.config/autostart/setxkbmap.desktop
[Desktop Entry]
Exec=/usr/bin/setxkbmap -option compose:ralt -option caps:escape_shifted_capslock
Name=setxkbmap configuration
Type=Application
